# supervisord.worker.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[program:celery.worker]
command=bash -c 'celery -A app.celery.worker.celery worker --loglevel=info --autoscale=%(ENV_CELERY_CONCURRENCY)s,1'
directory=/app
priority=2
autostart=true
autorestart=true
stdout_logfile=/var/log/worker_celery.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/worker_celery-error.log
stderr_logfile_maxbytes=50MB
environment=PYTHONUNBUFFERED=1